# crontab
0 * * * * /usr/bin/python3 /app/macos_version_checker.py
